# linux/system/bore-tunnel.service
[Unit]
Description=Pocket Lab — bore SSH tunnel
Documentation=https://github.com/Tsukieomie/pocket-lab-v2-6
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Read config from ~/.bore_env (BORE_HOST, BORE_SECRET override defaults)
EnvironmentFile=-%h/.bore_env
Environment=BORE_HOST=188.93.146.98
Environment=BORE_SECRET=pocketlab2026
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

ExecStart=/bin/sh -c '\
  REPO="%h/pocket-lab-v2-6"; \
  BORE_BIN=""; \
  for P in "$REPO/bore-custom-8443" "$REPO/bore-custom-2222" "%h/.local/bin/bore" "/usr/local/bin/bore" "/usr/bin/bore"; do \
    [ -x "$P" ] && BORE_BIN="$P" && break; \
  done; \
  if [ -z "$BORE_BIN" ]; then \
    echo "[bore-tunnel] bore not found — run: bash $REPO/linux/tunnel.sh install-bore"; \
    exit 1; \
  fi; \
  SECRET_ARG=""; [ -n "${BORE_SECRET}" ] && SECRET_ARG="--secret ${BORE_SECRET}"; \
  HOST_IP=$(getent ahostsv4 "${BORE_HOST}" 2>/dev/null | awk "/STREAM/{print \$1; exit}" \
    || python3 -c "import socket; print(socket.getaddrinfo(\"${BORE_HOST}\",None,socket.AF_INET)[0][4][0])" 2>/dev/null \
    || echo "${BORE_HOST}"); \
  echo "[bore-tunnel] Using binary: $BORE_BIN"; \
  echo "[bore-tunnel] Connecting to ${BORE_HOST} ($HOST_IP)"; \
  exec "$BORE_BIN" local 22 --to "$HOST_IP" $SECRET_ARG'

# Poll up to 15s for "listening at" line, then push bore-port.txt
ExecStartPost=/bin/sh -c '\
  TUNNEL_SH="%h/pocket-lab-v2-6/linux/tunnel.sh"; \
  LIVE_PORT=""; \
  for i in $(seq 1 15); do \
    sleep 1; \
    LIVE_PORT=$(journalctl --user -u bore-tunnel.service -n 200 --no-pager 2>/dev/null \
      | sed "s/\x1b\[[0-9;]*m//g" \
      | grep "listening at" | grep -oE ":[0-9]+" | tail -1 | tr -d ":" | grep -E "^[0-9]+$" \
      || echo ""); \
    [ -n "$LIVE_PORT" ] && break; \
  done; \
  if [ -n "$LIVE_PORT" ]; then \
    echo "[bore-tunnel] Port ready: $LIVE_PORT — syncing..."; \
    bash "$TUNNEL_SH" sync-port 2>&1 || true; \
  else \
    echo "[bore-tunnel] WARNING: no port reported within 15s"; \
  fi'

# Mark bore-port.txt as DOWN on stop
ExecStopPost=/bin/sh -c '\
  PORT_FILE="%h/pocket-lab-v2-6/bore-port.txt"; \
  [ -f "$PORT_FILE" ] || exit 0; \
  sed -i "s/^port=.*/port=DOWN/" "$PORT_FILE"; \
  sed -i "s/^updated=.*/updated=$(date -u +%Y-%m-%dT%H:%M:%SZ)/" "$PORT_FILE"; \
  echo "[bore-tunnel] bore-port.txt marked DOWN"'

# Always restart — tunnel must stay up
Restart=always
RestartSec=5
RestartMaxDelaySec=30

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
