# linux/system/bore-port-watcher.service
[Unit]
Description=Pocket Lab — bore port watcher (pushes port changes to GitHub)
Documentation=https://github.com/Tsukieomie/pocket-lab-v2-6
After=bore-tunnel.service
BindsTo=bore-tunnel.service

[Service]
Type=simple

EnvironmentFile=-%h/.bore_env
Environment=BORE_HOST=bore.pub
Environment=REPO_DIR=%h/pocket-lab-v2-6
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

# Poll bore-tunnel journal every 2s.
# On port change: call sync-port, retry up to 3 times with 5s backoff.
# On port loss: reset state so the next reconnect is detected cleanly.
ExecStart=/bin/sh -c '\
  TUNNEL_SH="${REPO_DIR}/linux/tunnel.sh"; \
  LAST_PORT=""; \
  echo "[watcher] started — watching bore-tunnel.service for port changes"; \
  while true; do \
    LIVE=$(journalctl --user -u bore-tunnel.service -n 200 --no-pager 2>/dev/null \
      | sed "s/\x1b\[[0-9;]*m//g" \
      | grep -oE "bore\.pub:[0-9]+" | tail -1 | cut -d: -f2 | grep -E "^[0-9]+$" \
      || journalctl --user -u bore-tunnel.service -n 200 --no-pager 2>/dev/null \
      | sed "s/\x1b\[[0-9;]*m//g" \
      | grep -oE "remote_port=[0-9]+" | tail -1 | cut -d= -f2 | grep -E "^[0-9]+$" \
      || echo ""); \
    if [ -n "$LIVE" ] && [ "$LIVE" != "$LAST_PORT" ]; then \
      echo "[watcher] port change: ${LAST_PORT:-none} → $LIVE"; \
      PUSH_OK=false; \
      for _try in 1 2 3; do \
        bash "$TUNNEL_SH" sync-port 2>&1 && PUSH_OK=true && break; \
        echo "[watcher] sync-port attempt $_try failed — retrying in 5s..."; \
        sleep 5; \
      done; \
      [ "$PUSH_OK" = true ] \
        && echo "[watcher] port $LIVE synced OK" \
        || echo "[watcher] WARNING: sync-port failed after 3 attempts (bore-port.txt local only)"; \
      LAST_PORT="$LIVE"; \
    elif [ -z "$LIVE" ] && [ -n "$LAST_PORT" ]; then \
      echo "[watcher] port lost — tunnel reconnecting"; \
      LAST_PORT=""; \
    fi; \
    sleep 2; \
  done'

Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=default.target
