# linux/system/cloudflared-tunnel.service
# linux/system/cloudflared-tunnel.service — Pocket Lab SSH Tunnel (cloudflared)
# systemd user unit — install to ~/.config/systemd/user/
#
# Uses Cloudflare Quick Tunnels over HTTPS/443.
# No account needed. Bypasses ISP port blocks (bore 7835 blocked, etc.)
#
# Install:
#   cp ~/pocket-lab-v2-6/linux/system/cloudflared-tunnel.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now cloudflared-tunnel.service
#
# Remove old bore-tunnel if previously installed:
#   systemctl --user disable --now bore-tunnel.service
#
# Logs:
#   journalctl --user -u cloudflared-tunnel.service -f
#
# After tunnel is up, run:
#   bash ~/pocket-lab-v2-6/linux/tunnel.sh sync-port
# to push the .trycloudflare.com hostname to bore-port.txt on GitHub.

[Unit]
Description=Pocket Lab — cloudflared SSH tunnel
Documentation=https://github.com/Tsukieomie/pocket-lab-v2-6
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Inherit GH_TOKEN, FS_BRIDGE_TOKEN, etc. from ~/.bore_env
EnvironmentFile=-%h/.bore_env
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

ExecStart=/bin/sh -c '\
  CF_BIN=""; \
  for P in "%h/.local/bin/cloudflared" "/usr/local/bin/cloudflared" "/usr/bin/cloudflared"; do \
    [ -x "$P" ] && CF_BIN="$P" && break; \
  done; \
  if [ -z "$CF_BIN" ]; then \
    echo "[cloudflared-tunnel] cloudflared not found — run:"; \
    echo "  bash %h/pocket-lab-v2-6/linux/tunnel.sh install-cloudflared"; \
    exit 1; \
  fi; \
  echo "[cloudflared-tunnel] Using binary: $CF_BIN"; \
  exec "$CF_BIN" tunnel --url tcp://localhost:22'

# After cloudflared reports a .trycloudflare.com hostname, push bore-port.txt
ExecStartPost=/bin/sh -c '\
  TUNNEL_SH="%h/pocket-lab-v2-6/linux/tunnel.sh"; \
  LIVE_HOST=""; \
  for i in $(seq 1 30); do \
    sleep 1; \
    LIVE_HOST=$(journalctl --user -u cloudflared-tunnel.service -n 300 --no-pager 2>/dev/null \
      | sed "s/\x1b\[[0-9;]*m//g" \
      | grep -oE "[a-z0-9-]+\.trycloudflare\.com" | tail -1 || echo ""); \
    [ -n "$LIVE_HOST" ] && break; \
  done; \
  if [ -n "$LIVE_HOST" ]; then \
    echo "[cloudflared-tunnel] Hostname ready: $LIVE_HOST — syncing..."; \
    bash "$TUNNEL_SH" sync-port 2>&1 || true; \
  else \
    echo "[cloudflared-tunnel] WARNING: no hostname reported within 30s"; \
  fi'

# Mark bore-port.txt as DOWN on stop
ExecStopPost=/bin/sh -c '\
  PORT_FILE="%h/pocket-lab-v2-6/bore-port.txt"; \
  [ -f "$PORT_FILE" ] || exit 0; \
  sed -i "s/^port=.*/port=DOWN/" "$PORT_FILE"; \
  sed -i "s/^host=.*/host=DOWN/" "$PORT_FILE"; \
  sed -i "s/^updated=.*/updated=$(date -u +%Y-%m-%dT%H:%M:%SZ)/" "$PORT_FILE"; \
  echo "[cloudflared-tunnel] bore-port.txt marked DOWN"'

# Always restart on failure — tunnel must stay up
Restart=always
RestartSec=5
RestartMaxDelaySec=60

# Safety
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
